# systemd/webcam-service.service
[Unit]
Description=WebSocket JSON-RPC 2.0 Server
Documentation=file:///opt/webcam-env/webcam_ip/README.md
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data

# Working directory
WorkingDirectory=/opt/webcam-env/webcam_ip

# Python executable and script
ExecStart=/opt/webcam-env/bin/python3 -m webcam_ip.server.websocket_server

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Environment variables
Environment=PYTHONPATH=/opt/webcam-env/webcam_ip
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/webcam-env/logs
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
RestrictNamespaces=true

# Process settings
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=websocket-jsonrpc

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
